# config/resfes.service
[Unit]
Description=Resfes service

[Service]
User=azureuser
ExecStart=/home/azureuser/resfes/start.sh

[Install]
WantedBy=multi-user.target

